# provision/ansible/roles/common/files/logrotate.cron
#!/bin/sh
logrotate /etc/logrotate.conf
EXITVALUE="${?}"

if [ X"${EXITVALUE}" != X"0" ]; then
    logger -t logrotate "ALERT exited abnormally with [${EXITVALUE}]"
fi
